# modules/adblock/adblock.cron
# HOMESERVER Adblock Update Cron
# Runs daily at 3 AM with random delay
0 3 * * * root sleep $((RANDOM % 3600)) && /usr/bin/python3 /var/local/lib/updates/modules/adblock/index.py >> /var/log/homeserver/adblock.log 2>&1 